# etc/cron.d/pkktime

SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin
MAILTO=root

 */5  *  *  *  *  root pkktime
